# config/com.logpose.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.logpose.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/remote-dev/log-pose/server/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/remote-dev/log-pose/server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>3000</string>
        <key>HOST</key>
        <string>127.0.0.1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/remote-dev/log-pose/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/remote-dev/log-pose/logs/server.error.log</string>

    <key>UserName</key>
    <string>remote-dev</string>

    <key>GroupName</key>
    <string>staff</string>
</dict>
</plist>
